# Universal-Binaries/WebDriver-387.10.10.10.40.140/Library/PreferencePanes/NVIDIA Driver Manager.prefPane/Contents/Resources/com.nvidia.nvagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.nvidia.nvagent</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/PreferencePanes/NVIDIA Driver Manager.prefPane/Contents/MacOS/NVAgent.app/Contents/MacOS/NVAgent</string>
	</array>
	<key>StartInterval</key>
	<integer>86400</integer>
	<key>QueueDirectories</key>
	<array>
		<string>/var/tmp/com.nvidia/com.nvidia.nvmenu</string>
		<string>/var/tmp/com.nvidia/com.nvidia.nvprefpane</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
