# Universal-Binaries/WebDriver-387.10.10.10.40.140/Library/LaunchDaemons/com.nvidia.nvroothelper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.nvidia.nvroothelper</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/PreferencePanes/NVIDIA Driver Manager.prefPane/Contents/MacOS/com.nvidia.nvroothelper</string>
	</array>
	<key>QueueDirectories</key>
	<array>
		<string>/var/tmp/com.nvidia/com.nvidia.nvagent</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
